# Linux-Fundamentals/crontab


crontab -e

* * * * * apt update -y > /var/log/syslog






